# rootfs/etc/services.d/PanelSense/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Example
# Runs python sense.py script
# ==============================================================================

bashio::log.info "Starting PanelSense..."

exec python3 /usr/src/app/senseapp/sense.py
