# panelsense/rootfs/etc/s6-overlay/s6-rc.d/init-panelsense/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: PanelSense
# Runs the Nginx daemon
# ==============================================================================

bashio::log.info "Panel Sense setup ..."
mkdir /config/panelsense

export PANEL_SENSE_DATABASE="/config/panelsense/sense_database.db"
